# effectiveMobile/monitoring_test.service
[Unit]
Description=Process Monitoring Service for "test"

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitoring_test.sh
